# ansible/roles/web/files/ncg-devops.ini
[program:ncg-devops]
directory=/var/www/ncg-devops
command=/var/www/ncg-devops/venv/bin/gunicorn --bind 0.0.0.0:8000 "flaskr:create_app()"
autostart=true
autorestart=true
stderr_logfile=/var/log/ncg-devops/ncg-devops.err.log
stdout_logfile=/var/log/ncg-devops/ncg-devops.out.log